# deploy/supervisor-agent.ini
[program:cloudflare-dns-agent]
directory=/www/wwwroot/cloudflare_dns/agent
command=/www/wwwroot/cloudflare_dns/.venv/bin/python agent.py
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=/www/wwwlogs/cloudflare-dns-agent.log
stderr_logfile=/www/wwwlogs/cloudflare-dns-agent-error.log
environment=CONTROL_URL="https://your-controller.example.com",AGENT_TOKEN="replace-me",PYTHONUNBUFFERED="1"

