# deploy/supervisor-backend.ini
[program:cloudflare-dns-backend]
directory=/www/wwwroot/cloudflare_dns/backend
command=/www/wwwroot/cloudflare_dns/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=/www/wwwlogs/cloudflare-dns-backend.log
stderr_logfile=/www/wwwlogs/cloudflare-dns-backend-error.log
environment=PYTHONUNBUFFERED="1"

